# escalation-league-backend/scripts/scryfall-update.timer
[Unit]
Description=Nightly Scryfall database update timer
Requires=scryfall-update.service

[Timer]
# Run at 6:00 AM UTC (1:00 AM EST) - after Scryfall's daily update
OnCalendar=daily
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
